# supervisord.conf
[program:builder]
command=/start web
directory=/app
autostart=true
autorestart=true
numprocs=1

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
numprocs=1
